# deploy/tencent/waronmaps.service
[Unit]
Description=War On Maps Rust game server
After=network.target

[Service]
Type=simple
User=__USER__
WorkingDirectory=__ROOT_DIR__
ExecStart=__ROOT_DIR__/rust_server/target/release/node_game_server 8002
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
